# etc/init/simplest.conf
description "The worlds' simplest webapp"
author "aurynn shaw"
start on runlevel [2345]
stop on shutdown

script
  exec /usr/bin/python /var/app/app.py
end script